# deployment/systemd/owlsome-frontend.service.example
[Unit]
Description=Owlsome Learning Frontend Preview
After=network.target

[Service]
Type=simple
WorkingDirectory=/data/workspace/projects/Owlsome/learning_platform/frontend
ExecStart=/usr/bin/npm run preview -- --host 0.0.0.0 --port 5173
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
